# aarch64-alpine-3.22.Dockerfile.d/200_seed_config.sh.d/etc/profile.d/umask.sh
# - /etc/profile.d/umask.sh -
umask 022
