# images/rpi64-rpios-lite-ci.Dockerfile.d/100_add_files.sh.d/etc/systemd/system/wpa_supplicant@wlan0.service.d/+force-interface-down.conf
[Service]
ExecPostDown=/sbin/ip link set down dev %I
