# images/rpi64-rpios-lite-ci.Dockerfile.d/100_add_files.sh.d/etc/systemd/system/wpa_supplicant@wlan0.service.d/+no-power-mgmt.conf
[Service]
ExecStartPre=/sbin/iw dev %i set power_save off             
