# turingpi2-alpine.Dockerfile.d/100_install_bmcd.sh.d/etc/profile.d/tpi.conf.sh
# this sets defaults for bmcd access via tpi
# enable and keep in sync with /etc/bmcd/config.yaml
# if changed

#export TPI_HOST=[::1]
#export TPI_PORT=443
