# rpi64-rpios-base.Dockerfile.d/100_add_files.sh.d/etc/systemd/system/wpa_supplicant@wlan0.service.d/+unblock-rfkill.conf
[Service]
ExecStartPre=/sbin/rfkill unblock wlan
