# amd64-EFI-CI-trixie.Dockerfile.d/100_add_files.sh.d/etc/systemd/system/getty@tty1.service.d/+ansible+systemd-consoles.conf
[Unit]
Description=Auto-Login root on %I

[Service]
ExecStart=
ExecStart=-/sbin/agetty -w -p -a root 38400 %I linux
