# turingpi2-alpine.Dockerfile.d/200_seed_config.sh.d/etc/init.d/coldplug
#!/sbin/openrc-run

description="Cold plug service"

depend()
{
  need dev
  keyword -docker -podman -lxc -openvz -prefix -systemd-nspawn -vserver
}

start()
{
  ebegin "Cold plug"
  for i in `find /sys -name "uevent"`; do
    ( 
      echo add > $i 
    ) 2> /dev/null
  done
  eend $?
  return 0
}
