# turingpi2-alpine.Dockerfile.d/200_seed_config.sh.d/etc/init.d/klogd
#!/sbin/openrc-run

command="/bin/busybox"
command_args="klogd"

depend() {
        need hostname syslog
        use clock
        before net
}
