# turingpi2-alpine.Dockerfile.d/200_seed_config.sh.d/etc/init.d/syslog
#!/sbin/openrc-run

# openrc for busysbox syslog
# it logs
#  * into local ring buffer and
#    This can be read via `busybox logger`
#  * to a syslog server on localhost for forwarding
#    TODO: for remote logging only enbale and configure rsyslog
#          rsyslog does also a bit syslog message format sanitisation

command="/bin/busybox syslogd"
command_args="-L -R 127.0.0.1:514 -C128 -S -t"

depend() {
        need hostname
        use clock
        before net
        provide logger
}
