# turingpi2-alpine.Dockerfile.d/200_seed_config.sh.d/etc/profile.d/aliases.sh
# - /etc/profile.d/aliases.sh -
alias rL=". /etc/profile"

#alias rS="exec screen -D -RR"
alias rS="exec tx"

#alias log="tail -n 99 -F /var/log/messages"
alias log="logread -F"
