# turingpi2-alpine.Dockerfile.d/200_seed_config.sh.d/etc/profile.d/aliases.sh
# - /etc/profile.d/aliases.sh -

alias rL=". /etc/profile"
alias rS="exec tx"
alias log="logread -F"
alias Faux='ps faux | grep [^]]$'
