# turingpi2-alpine.Dockerfile.d/200_seed_config.sh.d/etc/inittab
# /etc/inittab

::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default

::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/openrc shutdown

# allow password less logins on serial console
# as direct access to HW is already a security breach
ttyS0::respawn:/sbin/getty -L -w -n -l /lib/serial-login.sh ttyS0 115200 vt100
# USB-OTG
ttyGS0::respawn:/sbin/getty -L ttyGS0 115200 vt100
