# images/turingpi2-alpine.Dockerfile.d/100_install_bmcd.sh.d/etc/init.d/bmcd
#!/sbin/openrc-run

name="${RC_SVCNAME}"
command="/opt/bmcd/bin/bmcd"
supervisor=supervise-daemon

CERTFILE="/etc/ssl/certs/bmcd_cert.pem"
KEYFILE="/etc/ssl/certs/bmcd_key.pem"
CONFIG="/etc/bmcd/config.yaml"

command_args="--config $CONFIG"
command_background=true
pidfile="${RC_PREFIX}/run/bmcd.pid"

depend() {
  #need net # otg
  use logger
  after firewall
}

start_pre() {
  if [ ! -f "$CERTFILE" ] || [ ! -f "$KEYFILE" ]; then
    sh /etc/bmcd/generate_self_signedx509.sh
  fi
}

start_post() {
  # wait for bmcd listening, and execute startup config
  if [ -x /etc/bmcd/startup.sh ]; then
    for _ in 1 2 3 4 5 6 7 8 9 10; do
      if netstat -lntp | grep bmcd$ > /dev/null; then
        /etc/bmcd/startup.sh || :
        break
      fi
      sleep .5
    done
  fi
}
