# rpi64-raspios-lite-base.Dockerfile.d/100_add_files.sh.d/etc/systemd/system/ssh.service.d/ssh-keygen-A.conf
[Service]
ExecStartPre=
ExecStartPre=/bin/ssh-keygen -A
ExecStartPre=/sbin/sshd -t
