# bpi2m0-alpine.Dockerfile.d/200_seed_config.sh.d/etc/init.d/mmc
#!/sbin/openrc-run

description="mount MMC data partition"

depend() {
  need localmount dev
}

FILTER=${2:-LABEL=mmc-data}
LABEL="${FILTER#LABEL=}"

start() {
  ebegin "searching for data partition on MMC/SD - $FILTER"
  eval "$(blkid -t $FILTER -o export)"
  rc=$?
  eend $rc

  ebegin "mount $DEVNAME to /mmc"
  [ -d /mmc ] || \
    mkdir -p /mmc
  mount -t $TYPE "$DEVNAME" /mmc
  rc=$?
  eend $rc

  return $rc
}

stop() {
  ebegin "unmount /mmc"
  while umount /mmc 2> /dev/null; do : ; done
  grep " /mmc " /proc/mounts
  rc=$?
  case "$rc" in
    0 ) eend 1 ;;
    * ) eend 0 ;;
  esac
}

# vim: ts=2 et sw=2 ft=sh
