# alpine-base.Dockerfile.d/100_add_files.sh.d/etc/inittab
# /etc/inittab

::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default

# Set up a couple of gettys
# - no local gettys
#tty1::respawn:/sbin/getty 38400 tty1
#[...]
#tty6::respawn:/sbin/getty 38400 tty6

# Put a getty on the serial port
ttyS0::askfirst:/lib/serial-login.sh
#ttyGS0::respawn:/sbin/getty -L ttyGS0 115200 vt100

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
::shutdown:/sbin/openrc shutdown

