# docs/TuringPi2-Alpine/RK1-Armbian/base.patch.d/etc/profile.d/alias:ip=ip-6.sh
ip() (
  local opt=
  local af= # -6
  case "$1:$2" in
    a*: ) af= ; opt="-br" ;;
    -4:* | -6:* ) af="$1"; shift ;;
  esac
  exec /sbin/ip $af $opt "$@"
)
