# recipes/alpine/SBC-base/200_SBC-base.seed-config.sh.d/etc/profile.d/usr-local.sh
# - /etc/profile.d/usr-local.sh -
for i in \
  /usr/local/sbin \
  /usr/local/bin; \
do
  [ -d "$i" ] || continue
  export PATH="$i:$PATH"
done
