# images/rpi64-rpios-base.Dockerfile.d/100_add_files.sh.d/etc/systemd/system/etc-machine-id.service
[Unit]
Description=Ensure machine-id initialized
Before=systemd-networkd.service
After=local-fs.target
ConditionPathIsReadWrite=/etc
ConditionPathExists=!/etc/machine-id

[Install]
WantedBy=sysinit.target

[Service]
RemainAfterExit=no
Type=oneshot
ExecStart=/bin/systemd-machine-id-setup
