# images/bpi2m0-armbian-bookworm-base.Dockerfile.d/100_add_files.sh.d/etc/systemd/system/serial-getty@ttyS0.service.d/+ansible+systemd-consoles.conf
[Unit]
Description=Auto-Login root on %I

[Service]
ExecStart=
ExecStart=-/sbin/agetty -w -p -a root -s 115200,38400,9600 %I $TERM

